# packages/odoopbx/odoopbx-0.45.tar.gz/odoopbx-0.45/salt/roots/odoo/templates/odoo.service
[Unit]
Description=Odoo Service
After=network.target

[Service]
Type=simple
User={{ user }}
ExecStart={{ path }}/odoo/odoo-bin --config={{ path }}/odoo.conf --syslog
Restart=on-failure
RestartSec=10
SyslogIdentifier={{ odoover }}
StandardError=journal
StandardOutput=journal

[Install]
WantedBy=multi-user.target
