# packages/botlib/botlib-88.tar.gz/botlib-88/botd.service
[Unit]
Description=BOTD - the 24/7 channel daemon
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/bot mods=irc,rss,cmd

[Install]
WantedBy=multi-user.target
